# deploy/alphadesk.service
[Unit]
Description=AlphaDesk Quant Trading System
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=YOUR_VPS_USER
WorkingDirectory=/home/YOUR_VPS_USER/alphadesk
ExecStart=/home/YOUR_VPS_USER/alphadesk/venv/bin/python scheduler.py

# Environment variables (or use .env file)
EnvironmentFile=/home/YOUR_VPS_USER/alphadesk/.env

# Auto-restart on failure
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alphadesk

# Resource limits
MemoryMax=1G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
